# ops/systemd/masso-update.path
[Unit]
Description=Watch for masso-web update trigger

[Path]
PathExists=__REPO_ROOT__/ops/update/apply.trigger
Unit=masso-update.service

[Install]
WantedBy=multi-user.target
